# nimbus/confs/nimbus.cron
#!/bin/sh
su - postgres -c '/usr/bin/pg_dump bacula | bzip2 -9 > /var/nimbus/bacula-sql.bz2';
su - postgres -c '/usr/bin/pg_dump nimbus | bzip2 -9 > /var/nimbus/nimbus-sql.bz2';

/bin/chown nimbus.nimbus /var/nimbus/bacula-sql.bz2;
/bin/chown nimbus.nimbus /var/nimbus/nimbus-sql.bz2;

su - nimbus -c '/var/www/nimbus --update-graphs-data';
su - nimbus -c '/var/www/nimbus --delete-volumes';

