# nimbus/confs/nimbus.cron
#!/bin/sh
su - nimbus -c '/usr/bin/mysqldump -uroot -pn1mbus bacula | bzip2 -9 > /var/nimbus/bacula-sql.bz2';
su - nimbus -c '/usr/bin/mysqldump -uroot -pn1mbus nimbus | bzip2 -9 > /var/nimbus/nimbus-sql.bz2';
su - nimbus -c '/var/www/nimbus --update-graphs-data';
su - nimbus -c '/var/www/nimbus --delete-volumes';

