# nimbus/confs/remotestorage.cron
#!/bin/sh
su - nimbus -c '/var/www/nimbus --disk-usage-alert';

