# deb/DEBIAN/prerm
#!/bin/bash

/usr/sbin/userdel nimbus;
/usr/sbin/update-rc.d -f nimbusmanager remove;
/usr/bin/truecrypt -t -d -f;
/etc/init.d/nimbusmanager stop || echo;
